# supervisord.conf
[program:light-box]
command = gunicorn server:app -b localhost:8000
directory = /root/raspberry-light-box
user = www-data