# supervisord.conf
[program:light-box]
command=python server.py