# supervisord.conf
[program:light-box]
command = python server.py
directory = /var/www/raspberry-light-box
user = www-data